# supervisord.conf
[supervisord]
nodaemon=true

[program:reset_tables]
command=python -u test.py
directory=/gesture_recognition/offline